# services/evolution-source-registry/evolution-source-registry.service.in
[Unit]
Description=Evolution source registry

[Service]
Type=dbus
BusName=@SOURCES_DBUS_SERVICE_NAME@
ExecStart=@libexecdir@/evolution-source-registry
