# services/evolution-calendar-factory/evolution-calendar-factory.service.in
[Unit]
Description=Evolution calendar service

[Service]
Type=dbus
BusName=@CALENDAR_DBUS_SERVICE_NAME@
ExecStart=@libexecdir@/evolution-calendar-factory
