# services/evolution-user-prompter/evolution-user-prompter.service.in
[Unit]
Description=Evolution user prompter

[Service]
Type=dbus
BusName=@USER_PROMPTER_DBUS_SERVICE_NAME@
ExecStart=@libexecdir@/evolution-user-prompter
